# roles/podman_nginx/templates/certbot-renew.timer
[Unit]
Description=Timer for certbot renewals

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
